# .config/autostart/blueman-applet.desktop
[Desktop Entry]
Exec=blueman-applet
Name=blueman-applet
Type=Application
