# .config/autostart/dbus-update-activation-environment.desktop
[Desktop Entry]
Exec=dbus-update-activation-environment --all
Name=dbus-update-activation-environment
Type=Application
