# .config/autostart/emacs.desktop
[Desktop Entry]
Exec=emacs --daemon
Name=emacs
Comment=GNU project Emacs editor
Type=Application
