# .config/autostart/gammastep-indicator.desktop
[Desktop Entry]
Exec=gammastep-indicator
Name=gammastep-indicator
Type=Application
