# .config/autostart/gnome-keyring-daemon.desktop
[Desktop Entry]
Exec=gnome-keyring-daemon
Name=gnome-keyring-daemon
Type=Application
