# .config/autostart/kanshi.desktop
[Desktop Entry]
Exec=kanshi
Name=kanshi
Type=Application
