# .config/autostart/lxqt-policykit-agent.desktop
[Desktop Entry]
Exec=lxqt-policykit-agent
Name=lxqt-policykit-agent
Comment=Polkit agent for the LXQt desktop environment
Type=Application
