# .config/autostart/nextcloud.desktop
[Desktop Entry]
Exec=nextcloud --background
Name=nextcloud
Type=Application
