# .config/autostart/solaar.desktop
[Desktop Entry]
Exec=solaar
Name=solaar
Type=Application
