# .config/autostart/swaync.desktop
[Desktop Entry]
Exec=swaync
Name=swaync
Type=Application
