# .config/autostart/swaybg.desktop
[Desktop Entry]
Exec=swaybg -i /home/flake/.background-image
Name=swaybg
Type=Application
