# .config/autostart/udiskie.desktop
[Desktop Entry]
Exec=udiskie --appindicator -t
Name=udiskie
Type=Application
