# .config/autostart/udiskie.desktop
[Desktop Entry]
Exec=udiskie
Name=udiskie
Type=Application
