# .config/autostart/waybar.desktop
[Desktop Entry]
Exec=waybar
Name=waybar
Type=Application
