# .config/autostart/wl-paste.desktop
[Desktop Entry]
Exec=wl-paste --watch cliphist store
Name=wl-paste
Type=Application
